# skills/devops/cloudflare-tunnel-setup/cloudflared-trading.service
[Unit]
Description=Cloudflare Tunnel (trading)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/terexitarius/.cloudflared/config.yml run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
